# deploy/leakx-webhooks.service
[Unit]
Description=LeakX - livraison des webhooks
After=network.target

[Service]
Type=oneshot
User=leakx
Group=leakx
WorkingDirectory=/opt/leakx/api
EnvironmentFile=/opt/leakx/api/.env
ExecStart=/opt/leakx/api/.venv/bin/python scripts/deliver_webhooks.py
